# build/supervisor/supervisor.conf
[program:scheduler]
process_name=%(program_name)s_%(process_num)02d
directory=/app
command=/bin/sh -c "while [ true ]; do (php artisan schedule:run --verbose --no-interaction); sleep 60; done"
autostart=true
autorestart=true
numprocs=1
user=kdock
redirect_stderr=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout

[program:default]
process_name=%(program_name)s_%(process_num)02d
directory=/app
command=php artisan queue:work redis --sleep=3 --tries=1
autostart=true
autorestart=true
numprocs=2
user=kdock
redirect_stderr=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout

[program:og_api_jobs]
process_name=%(program_name)s_%(process_num)02d
directory=/app
command=php artisan queue:work --queue=og_api_jobs --sleep=3 --tries=1
autostart=true
autorestart=true
numprocs=1
user=kdock
redirect_stderr=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout

[program:db_jobs]
process_name=%(program_name)s_%(process_num)02d
directory=/app
command=php artisan queue:work --queue=db_jobs --sleep=3 --tries=1
autostart=true
autorestart=true
numprocs=1
user=kdock
redirect_stderr=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout

[program:mail]
process_name=%(program_name)s_%(process_num)02d
directory=/app
command=php artisan queue:work --queue=mail --sleep=3 --tries=1
autostart=true
autorestart=true
numprocs=2
user=kdock
redirect_stderr=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout

[program:logs]
process_name=%(program_name)s_%(process_num)02d
directory=/app
command=php artisan queue:work --queue=logs --sleep=3 --tries=1
autostart=true
autorestart=true
numprocs=3
user=kdock
redirect_stderr=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout